# serve.service
[Unit]
Description=Serve - Static File Server
After=network.target

[Service]
Type=simple
User=serve
Group=serve
ExecStart=/usr/local/bin/serve -c /etc/serve/serve.toml
WorkingDirectory=/var/www/serve
Restart=on-failure
RestartSec=5

# Allow binding to privileged ports (443) as non-root
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Hardening
ProtectSystem=strict
ProtectHome=read-only
NoNewPrivileges=false
ReadWritePaths=/var/log/serve /var/www/serve
ReadOnlyPaths=/etc/letsencrypt

[Install]
WantedBy=multi-user.target
